# ks-manage/golden-boot-templates/golden-boot.service
[Unit]
Description=Golden Boot Initialization
After=network-online.target systemd-user-sessions.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash /root/golden-boot.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
