# ks-manage/addons-for-kickstarts/tux2lab-sync.timer
[Unit]
Description=tux2lab configuration sync timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
